# weather.timer
[Unit]
Description=Measure weather stats periodically

[Timer]
OnUnitActiveSec=1h
Persistent=true

[Install]
WantedBy=timers.target
